# root/etc/services.d/transmission/run
#!/usr/bin/with-contenv bash

exec \
      s6-setuidgid abc /usr/bin/transmission-daemon -f \
      --log-error -g /config --logfile /config/transmission.log
